# etc/init/aquaria.conf.in
# Aquaria
#
# This service maintains an Aquaria service from the point the system is
# started until it is shut down again.

start on stopped rc RUNLEVEL=[2345]
stop on runlevel [!2345]

respawn
exec @bindir@/aquaria
